# cron.d/qos-rrd_graph
*/5 * * * * root /srv/qos/bin/rrd-graph.sh > /dev/null
