# cron.d/qos-rrd_update
* * * * * root /srv/qos/bin/rrd-update.sh
